# bridge/com.albertying.tab-out-tree.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.albertying.tab-out-tree</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/A.Y/.nvm/versions/node/v22.22.0/bin/node</string>
        <string>/Users/A.Y/programs/tab-out/bridge/claude-bridge.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/A.Y/.local/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/A.Y</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/A.Y/Library/Logs/tab-out-tree.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/A.Y/Library/Logs/tab-out-tree.log</string>
</dict>
</plist>
